# charm/layer/turnip-storage/templates/data.mount.j2
[Unit]
Description=Turnip data file system

[Mount]
What={{ hostname }}:{{ mountpoint }}
Where={{ data_dir }}
Type={{ fstype }}
Options={{ options }}

